# sv/polkit-gnome/run
#!/bin/execline -P
exec /usr/libexec/polkit-gnome-authentication-agent-1
