# sv/mpdris2/run
#!/bin/sh
mpdris2-rs --host /run/user/1000/mpd.sock

